# MacOS/restic-backup-local-03.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>restic-backup-local-03</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/lsampaio/Luciano/git/home/backup-restic-s3-minio/restic-backup.sh</string>
      <string>local</string>
      <string>/Volumes/Backup-03/MacOS-Backup-Luciano</string>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/lsampaio/Luciano/git/home/backup-restic-s3-minio/logs/backup-local-03.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/lsampaio/Luciano/git/home/backup-restic-s3-minio/logs/backup-local-03.log</string>
    <key>RunAtLoad</key>
    <false />
    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
      <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
      <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
      <dict>
        <key>Hour</key>
        <integer>14</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
      <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
      <dict>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
    </array>
  </dict>
</plist>